# systemd/llama-cpp-server-misgastos.service
[Unit]
Description=llama.cpp server para MisGastos (Qwen3.5-9B)
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/llama-server \
 -m %h/.cache/llama.cpp/models/Qwen3.5-9B-Q4_K_M.gguf \
 --mmproj %h/.cache/llama.cpp/models/mmproj-BF16.gguf \
 --port 8005 \
 --host 0.0.0.0 \
 --ctx-size 16384 \
 --gpu-layers 99 \
 --threads 16 \
 --batch-size 512 \
 --ubatch-size 512 \
 --flash-attn on
Restart=on-failure
RestartSec=5
Environment=LLAMA_SERVER_LOG_LEVEL=info

[Install]
WantedBy=default.target
